# init.d/init.d.nginx
#!/bin/sh
### BEGIN INIT INFO
# Provides:          nginx
# Required-Start:    $network $remote_fs $local_fs 
# Required-Stop:     $network $remote_fs $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start and stop the custom compiled NGINX server
### END INIT INFO

# Author:   getlnmp
# website:  https://getlnmp.com

PATH=/sbin:/usr/sbin:/bin:/usr/bin

# --- PATHS FOR YOUR CUSTOM NGINX ---
NAME="nginx"
DESC="Nginx Web Server"
CONFFILE="/usr/local/nginx/conf/nginx.conf"
DAEMON="/usr/local/nginx/sbin/nginx"
PIDFILE="/usr/local/nginx/logs/nginx.pid"

# Default variables
SLEEPSEC=1
CHECKSLEEP=3
RETVAL=0

# Check if the daemon exists
[ -x $DAEMON ] || exit 0

# --- CORE FUNCTIONS ---

# Function to test the config file
do_configtest() {
    $DAEMON -t -c $CONFFILE
    RETVAL="$?"
    return $RETVAL
}

# Function to start the service
do_start()
{
    # Check 1: If the PID file exists and process is running
    if [ -f $PIDFILE ] && /bin/ps -p $(/bin/cat $PIDFILE) > /dev/null; then
        echo "$DESC is already running (PID $(/bin/cat $PIDFILE))."
        return 0 # Success: already running
    fi

    do_configtest || return 6 # Check config before starting

    echo -n "Starting $DESC: "
    # Run the daemon, specify the config file, and ensure it runs in the background (-g 'daemon on;')
    $DAEMON -c $CONFFILE
    RETVAL="$?"
    if [ "$RETVAL" = 0 ]; then
        echo "OK"
    else
        echo "FAILED"
    fi
    return $RETVAL
}

# Function to stop the service gracefully
do_stop()
{
    echo -n "Stopping $DESC: "
    
    # Check 1: If the PID file doesn't exist OR process is NOT running
    if ! [ -f $PIDFILE ] || ! /bin/ps -p $(/bin/cat $PIDFILE) > /dev/null; then
        echo "NOT RUNNING"
        /bin/rm -f $PIDFILE # Clean up any stale PID file
        return 0 # Success: already stopped
    fi

    # If we reach here, the service is running. Proceed with graceful stop.
    
    # Send SIGTERM (graceful stop)
    /bin/kill -s TERM $(/bin/cat $PIDFILE)
    
    # Wait for the process to exit
    COUNT=0
    while [ $COUNT -lt 15 ] && /bin/ps -p $(/bin/cat $PIDFILE) > /dev/null; do
        /bin/sleep 1
        COUNT=$((COUNT + 1))
    done

    if /bin/ps -p $(/bin/cat $PIDFILE) > /dev/null; then
        echo -n "FAILED (forcing kill) "
        # If still running, send SIGKILL
        /bin/kill -s KILL $(/bin/cat $PIDFILE)
        RETVAL=2
    else
        echo "OK"
        RETVAL=0
    fi
    /bin/rm -f $PIDFILE
    return $RETVAL
}

# Function to reload the configuration
do_reload() {
    do_configtest || return 6 # Check config before reloading
    echo -n "Reloading $DESC configuration: "
    if [ -f $PIDFILE ]; then
        # Send SIGHUP (graceful reload)
        /bin/kill -s HUP $(/bin/cat $PIDFILE)
        echo "OK"
        RETVAL=0
    else
        echo "NOT RUNNING (cannot reload)"
        RETVAL=1
    fi
    return $RETVAL
}

# --- CONTROL FLOW ---
case "$1" in
    start)
        do_start
        ;;
    stop)
        do_stop
        ;;
    status)
        # Check if the PID file exists and the process is running
        if [ -f $PIDFILE ] && /bin/ps -p $(/bin/cat $PIDFILE) > /dev/null; then
            echo "$NAME is running."
            exit 0
        else
            echo "$NAME is NOT running."
            exit 3
        fi
        ;;
    configtest)
        do_configtest
        ;;
    reload)
        do_reload
        ;;
    restart)
        do_stop
        sleep 1
        do_start
        ;;
    check-reload)
        do_checkreload
        RETVAL=0
        ;;
    *)
        echo "Usage: $NAME {start|stop|status|restart|reload|configtest}" >&2
        exit 3
        ;;
esac

exit $RETVAL